# systemd/otoparcapanel-monitor.service
[Unit]
Description=Oto Parça Panel - System Monitor Service
Requires=otoparcapanel.service
After=otoparcapanel.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/otoparcapanel

# Monitor script
ExecStart=/opt/otoparcapanel/scripts/monitor.sh

# Restart policy
Restart=always
RestartSec=30

# Timeout
TimeoutStartSec=60
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=otoparcapanel-monitor

# Environment
Environment=MONITOR_INTERVAL=60
Environment=ALERT_EMAIL=admin@example.com
Environment=DISK_THRESHOLD=85
Environment=MEMORY_THRESHOLD=90
Environment=CPU_THRESHOLD=95

[Install]
WantedBy=multi-user.target