# systemd/otoparcapanel-backup.service
[Unit]
Description=Oto Parça Panel - Database Backup Service
Requires=otoparcapanel.service
After=otoparcapanel.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/otoparcapanel

# Backup script
ExecStart=/opt/otoparcapanel/scripts/backup.sh

# Timeout
TimeoutStartSec=600

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=otoparcapanel-backup

# Environment
Environment=BACKUP_RETENTION_DAYS=30
Environment=BACKUP_DIR=/opt/otoparcapanel/backups